# crontabs
# test
* * * * * /bin/date > /var/log/date.log 2>&1 